# shellscript/startup/myshell.service
[Unit]
Type=simple
Description=my-script shell

[Service]
ExecStart=/bin/bash /data/starting.sh

[Install]
WantedBy=multi-user.target